# env/supervisor.conf
[program:web]
autorestart = true
priority = 10
numprocs = 1
user = ubuntu
directory = /home/ubuntu/service
command = /home/ubuntu/service/run.sh
environment=HOME="/home/ubuntu",USER="ubuntu"
